# files/system/etc/rc.d/rc.local
#!/bin/bash

exec bash /opt/zapret/init.d/sysv/zapret start
